# fr24/files/etc/network/if-pre-up.d/iptables
#!/bin/bash
/sbin/iptables-restore < /etc/network/iptables
/sbin/ip6tables-restore < /etc/network/ip6tables
